# jackett.conf
[program:jackett]
user = nobody
command = /usr/bin/mono --debug /opt/Jackett/Jackett.exe
